# scripts/init/gateone.conf
# Gate One - Web-based terminal and application platform

# NOTE: Gate One settings can be found in /etc/gateone/conf.d
# Documentation:  http://liftoff.github.io/GateOne/

description "Gate One"
author "Dan McDougall"

start on (starting network-interface
          or starting network-manager
          or starting networking)

stop on runlevel [!023456]

respawn
respawn limit 10 5
umask 022

exec gateone
